# .config/systemd/user/scr.service
[Unit]
Description=Persistent Scratchpad

[Service]
ExecStart=/home/aditya/.scripts/scr_exec
Restart=always

[Install]
WantedBy=multi-user.target
